# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid
user=root

[program:ttyd]
command=ttyd -p 7681 -W zsh
autostart=true
autorestart=true
startsecs=3
startretries=999
stderr_logfile=/data/ttyd.err.log
stdout_logfile=/data/ttyd.out.log
user=root
environment=
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",
    NPM_CONFIG_PREFIX="/usr/local/share/npm-global",
    PATH="/usr/local/share/npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",
    HOME="/root",
    SHELL="/bin/zsh"
```

---

## Railway Restart Policy

Railway Dashboard → your service → **Settings → Deploy**:

| Setting | Value |
|---|---|
| **Restart Policy** | `Always` |
| **Max Restarts** | `unlimited` or highest value |
| **Healthcheck Path** | leave blank |
| **Healthcheck Timeout** | `300` seconds |

---

## Complete Crash Protection Stack
```
ttyd crashes
    → supervisord restarts it instantly (startretries=999)
        → if whole container crashes
            → Railway restart policy brings it back
                → /data volume means nothing is lost
```

Three layers of protection. It will never stay down.

---

## File Structure in Your Repo
```
.devcontainer/
├── Dockerfile
└── supervisord.conf     ← new file to create
