# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
startretries=999
stderr_logfile=/data/nginx.err.log
stdout_logfile=/data/nginx.out.log

[program:ttyd]
command=ttyd -p 7681 -W tmux new-session -A -s main
autostart=true
autorestart=true
startsecs=3
startretries=999
stderr_logfile=/data/ttyd.err.log
stdout_logfile=/data/ttyd.out.log
user=root
environment=
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",
    NPM_CONFIG_PREFIX="/usr/local/share/npm-global",
    PATH="/usr/local/share/npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",
    HOME="/root",
    SHELL="/bin/zsh"

[program:filebrowser]
command=filebrowser -r /data/workspace -p 8081 -a 0.0.0.0 --noauth -d /data/filebrowser/filebrowser.db
autostart=true
autorestart=true
startsecs=3
startretries=999
stderr_logfile=/data/filebrowser.err.log
stdout_logfile=/data/filebrowser.out.log
user=root
nodaemon=true
logfile=/data/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
startretries=999
stderr_logfile=/data/nginx.err.log
stdout_logfile=/data/nginx.out.log

[program:ttyd]
command=ttyd -p 7681 -W tmux new-session -A -s main
autostart=true
autorestart=true
startsecs=3
startretries=999
stderr_logfile=/data/ttyd.err.log
stdout_logfile=/data/ttyd.out.log
user=root
environment=
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",
    NPM_CONFIG_PREFIX="/usr/local/share/npm-global",
    PATH="/usr/local/share/npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",
    HOME="/root",
    SHELL="/bin/zsh"

[program:filebrowser]
command=filebrowser -r /data/workspace -p 8081 -a 0.0.0.0 --noauth -d /data/filebrowser/filebrowser.db
autostart=true
autorestart=true
startsecs=3
startretries=999
stderr_logfile=/data/filebrowser.err.log
stdout_logfile=/data/filebrowser.out.log
user=root
